# spider_manager.conf
[program:spider_manager]
command=python3 -u /bots/weibo0525.1/manager.py
autorstart=true
autorestart=true
